# crontab
* * * * * . $HOME/.bash_profile ; DISPLAY=:0 $HOME/.bin/changeWP
